# templates/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:keystone]
user=keystone
command=/usr/bin/keystone-all
